# systemd/keepalive-client.timer
[Unit]
Description=Send keepalive info

[Timer]
OnCalendar=hourly
AccuracySec=1h

[Install]
WantedBy=timers.target
