# systemd/keepalive-client.timer
[Unit]
Description=Send keepalive info

[Timer]
OnCalendar=*:0,15,30,45
AccuracySec=1m

[Install]
WantedBy=timers.target
