# systemd/keeplive-client.timer
[Unit]
Description=Send keepalive info

[Timer]
OnCalendar=daily
RandomizedDelaySec=12m
AccuracySec=1h

[Install]
WantedBy=timers.target
